# docker/crontab
# Auto build from buildContainer script
PYTHONPATH=/app/appRequirements 
 
# Run automation x every min
*/3 * * * * cd /app && /app/runMain.sh
